# services/USERNAME_HERE/USERNAME_HERE_dev.service
[Unit]
Description=Travian game engine (dev)
After=network.target
After=mysqld.service

[Service]
User=travian
Type=forking
ExecStart=/home/travian/USERNAME_HERE/servers/dev/include/USERNAME_HERE_dev.service.php
Restart=always

[Install]
WantedBy=multi-user.target
